# board/rjg/copilot_port_tester/rootfs_overlay/etc/inittab
::sysinit:/etc/sysinit
::restart:/sbin/init
::shutdown:sync
::shutdown:umount -a -r
::ctrlaltdel:/sbin/reboot
::respawn:/etc/start_task
